# script/ansible/roles/tmsetup/templates/services/traffic-monitor.target.j2
[Unit]
Description=Traffic Monitor parent service profile

[Install]
WantedBy=go2rtc_server.service frigate.service node-red-tm.service plate-recognizer.service
